# conf/systemd.service
[Unit]
Description=Galène: Videoconferencing server
Documentation=https://galene.org
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__FINALPATH__/live/
ExecStart=__FINALPATH__/live/galene -http "127.0.0.1:__PORT__" -insecure -turn __PUBLIC_IP4__:__TURN_PORT__ -udp-range 49152-65535 -groups __DATADIR__/groups -recordings __DATADIR__/recordings  -data __FINALPATH__/live/data/
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
