# conf/systemd.service
[Unit]
Description=Galène: videoconferencing server
Documentation=https://galene.org
After=network.target

[Service]
Type=simple
WorkingDirectory=__FINALPATH__/
User=__APP__
Group=__APP__
ExecStart=__FINALPATH__/galene -turn __PUBLIC_IP4__.1:1194
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
