# conf/systemd.service
[Unit]
Description=Galene
After=network.target

[Service]
Type=simple
WorkingDirectory=__FINALPATH__/
User=__APP__
Group=__APP__
ExecStart=__FINALPATH__/galene
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target